# supervisord/init.conf
[program:init]
command=/init.sh
autorestart=false
redirect_stderr=true
redirect_stdout=/dev/stdout
