# supervisord/php-service.conf
[program:php]
command=/php.sh
autorestart=false
redirect_stderr=true
redirect_stdout=/dev/stdout
